# backend/deploy/qdrant.service
[Unit]
Description=Qdrant Vector Database
Documentation=https://qdrant.tech
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=qdrant
Group=qdrant
WorkingDirectory=/var/lib/qdrant
ExecStart=/usr/local/bin/qdrant --config-path /etc/qdrant/config.yaml
ExecReload=/bin/kill -HUP $MAINPID

# 重启策略
Restart=always
RestartSec=5
StartLimitBurst=5

# 输出重定向
StandardOutput=append:/var/log/qdrant/qdrant.log
StandardError=append:/var/log/qdrant/qdrant-error.log
SyslogIdentifier=qdrant

# 环境变量
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1
Environment=QDRANT_LOG_LEVEL=INFO

# 安全配置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/qdrant /var/log/qdrant

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0

# 内存限制（根据服务器配置调整）
# MemoryMax=4G

[Install]
WantedBy=multi-user.target 